# docker/app/queue/worker-default.conf
[program:worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/incul-agent/artisan queue:work --sleep=3 --tries=2 --max-time=3600 --timeout=3600 --queue=default
autostart=true
autorestart=true
numprocs=1
user=www
redirect_stderr=true
stdout_logfile=/var/www/html/incul-agent/storage/logs/queue-default.log

[program:scheduler]
process_name=%(program_name)s_%(process_num)02d
command=/bin/sh -c "while [ true ]; do (php /var/www/html/incul-agent/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
autostart=true
autorestart=true
numprocs=1
user=www
redirect_stderr=true
